# apps/filebrowser/root/etc/s6-overlay/s6-rc.d/kodebase/run
#!/command/with-contenv bash

ln -s /storage/config/plextraktsync /home/beepboop/.config

/kodebase.sh